# debian/debian/postinst
#!/bin/sh
update-rc.d anycast-dns-monitor defaults
