# debian/debian/postrm
#!/bin/sh
update-rc.d -f anycast-dns-monitor remove
